# resources/services/grafana.service
# container-5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613.service
# autogenerated by Podman 3.3.1
# Mon Oct  4 15:03:13 EDT 2021

[Unit]
Description=Podman container-5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/user/1000/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613
ExecStop=/usr/bin/podman stop -t 10 5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613
ExecStopPost=/usr/bin/podman stop -t 10 5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613
PIDFile=/run/user/1000/containers/overlay-containers/5db45a41e018c53171ea6a5a5b189dc3974a8d3eb8da1aad1dac0e657617b613/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
